# systemd/meetingbox-ui.service
[Unit]
Description=MeetingBox OLED Display UI
After=network.target meetingbox-backend.service
Requires=meetingbox-backend.service

[Service]
Type=simple
User=meetingbox
Environment="DISPLAY=:0"
Environment="BACKEND_URL=http://localhost:8000"
Environment="KIVY_NO_CONSOLELOG=1"
WorkingDirectory=/opt/meetingbox/device-ui
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 /opt/meetingbox/device-ui/src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meetingbox-ui

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
